# target/linux/tr6560/base-files/etc/hotplug.d/led/01-led_end
#!/bin/sh
. /usr/sbin/wps_led.sh

[ "$ACTION" = "led_end" ] || exit 0

logger -t led_end -p notice "ACTION=$ACTION" "led=$LED"

led off
