# target/linux/tr6560/generic/base-files/etc/uci-defaults/04_wireless
[ ! -e /etc/config/wireless ] && exit 0

. /lib/functions/system.sh



set_device_config()
{
    device=$(board_name)

    case $device in
        *)
            ;;
    esac
}



set_device_config


exit 0
