# target/linux/tr6560/base-files/etc/hotplug.d/working-mode/01_working-mode_ClearLuciCache
#!/bin/sh


# delete luci cache because agent shows different pages
rm /tmp/luci-* -rf

