# target/linux/tr6560/base-files/etc/hotplug.d/mesh/01_mesh_led
#!/bin/sh
. /usr/share/libubox/jshn.sh

logger -t mesh_led -p notice "ACTION=$ACTION"

mesh_led=`ubus call uci get '{"config":"system","section":"led_mesh","option":"sysfs"}' | jsonfilter -e "@.value"`

linkup()
{
    json_init
    json_add_object "leds"
    json_add_int "$mesh_led" 255
    json_close_object
    json_add_int "on" 0
    ubus call led set "$(json_dump)"
}

linkdown()
{
    json_init
    json_add_object "leds"
    json_add_int "$mesh_led" 0
    json_close_object
    ubus call led set "$(json_dump)"
}

linking()
{
    json_init
    json_add_object "leds"
    json_add_array "$mesh_led"
    json_add_int "" 0
    json_add_int "" 255
    json_close_array
    json_close_object
    json_add_int "on" 150
    json_add_int "off" 150
    json_add_int "blink" 120000
    ubus call led set "$(json_dump)"
}

case "$ACTION" in
    "linkup")
        linkup
    ;;
    "linkdown")
        linkdown
    ;;
    "linking")
        linking
    ;;
esac

