# target/linux/tr6560/base-files/etc/hotplug.d/mesh/02_mesh_role
#!/bin/sh
. /lib/functions/system.sh
. /lib/functions/working-mode.sh



case "$ACTION" in
    "rolechange")
        if [ "$MODE" != "gateway_mode" ]; then
            exit 0
        fi

        /bin/scripts/modeSet.sh  controller
        ;;
esac
